# rootfs/etc/services.d/nginx/run
#!/bin/sh

mkdir -p /run/nginx

exec nginx -g "daemon off;"
