# rootfs/etc/services.d/php-fpm/run
#!/bin/sh

exec php-fpm --nodaemonize 2>&1 | logger
